# cableos/scripts/cableos-install/cableos-installer.service
[Unit]
Description=CableOS Installation Single-Use Startup Script
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cableos-installer.sh -v -i
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
